# deploy/cloud-vps/systemd/wiki-econ-refresh.timer
[Unit]
Description=Nightly wiki-economics refresh

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
Unit=wiki-econ-refresh.service

[Install]
WantedBy=timers.target
